# final_project/main.cron
21 3 * * * /home/de11an/kart/project/main.py